# files/xinetd.d/vnc
service vnc
{
	disable		= no
	socket_type 	= stream
	protocol	= tcp
	wait		= no
	user		= root
	server		= /home/dsm/scripts/dummy.sh
	port		= 5900
	per_source	= 1
	instances	= 10
}
